# getsby-sync.1.timer
[Unit]
Description=getsby sync timer

[Timer]
OnBootSec=10min
OnUnitActiveSec=15min
Unit=getsby-sync.1.service
